# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket
[Socket]
ListenStream=/run/quora_clone-gunicorn.sock
[Install]
WantedBy=sockets.target